# contrib/iperf3_@.socket
[Socket]
ListenStream=%i

[Install]
WantedBy=sockets.target

