# contrib/iperf3.service
[Unit]
Description=iperf3
Requires=network-online.target

[Service]
ExecStart=/usr/bin/iperf3 -s

[Install]
WantedBy=multi-user.target
